# batch/rcd_server.in
#! /bin/bash
# armagetronad start script

# init info lifted from debian package rules
### BEGIN INIT INFO
# Provides:          @progname@
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs $network
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: @progtitle@ Dedicated Server
# Description:       Runs the @progtitle@ Dedicated Server
### END INIT INFO

# program installation directory
PROGDIR=@scriptdir_reloc@

# load configuration
. ${PROGDIR}/rcd_config

# files where PIDs are stored
STARTERPIDFILE=${PIDDIR}/@progname@-starter.pid
MAINPIDFILE=${PIDDIR}/@progname@.pid

# program to run
BINARY=serverstarter
BINFILE=${PROGDIR}/${BINARY}
TYPE=server

# log file
LOGFILE=${LOGDIR}/@prognamebase@.log

# delegate to common start/stop/status code
. ${RCDDIR}/rcd_startstop
